# packages/unregistry/output/DEBIAN/prerm
#!/bin/sh
set -e

case "$1" in
    remove|upgrade|deconfigure)
        systemctl stop unregistry.service || true
        systemctl disable unregistry.service || true
        systemctl daemon-reload
        ;;
    *)
        echo "prerm called with unknown argument \`$1'" >&2
        exit 1
        ;;
esac

exit 0
